# .config/systemd/user/systemd-gnome-ask-password-agent.path
[Unit]
Description=Forward Password Requests to systemd-gnome-ask-password-agent
DefaultDependencies=no
Conflicts=shutdown.target emergency.target

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes

[Install]
WantedBy=gnome-session-initialized.target
